# wscgo.service
[Unit]
Description=wscgo service

[Service]
Type=simple
ExecStart=/usr/bin/wscgo /etc/wscgo.ini

[Install]
WantedBy=multi-user.target
